# systemd/routeros-watchdog.service
[Unit]
Description=Smart Multi-WAN Router OS Health Monitor Watchdog
Documentation=https://github.com/yourusername/smart-multi-wan-router
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/routeros/watchdog
ExecStart=/usr/bin/python3 /opt/routeros/watchdog/watchdog_service.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=10

# Security settings
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /opt/routeros/config /var/run
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=routeros-watchdog

[Install]
WantedBy=multi-user.target